# lib/trundle.conf
#!upstart
description "node.js trundle process monitor"
author      "max"

start on startup
stop on shutdown

script
    export HOME="/home/max"
    TRUNDLE_DB="http://admin:pass@localhost:5984/processes"
    exec sudo -u max /usr/local/bin/trundle load $TRUNDLE_DB > /var/log/trundle.log
end script